# deploy/fretty.service
[Unit]
Description=Fretty - Guitar Fretboard Trainer
After=network.target

[Service]
User=fretty
Group=fretty
WorkingDirectory=/opt/fretty/backend
EnvironmentFile=/opt/fretty/.env
ExecStart=/opt/fretty/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8002
Restart=always
RestartSec=3

MemoryMax=512M
CPUQuota=80%
ProtectSystem=strict
ReadWritePaths=/opt/fretty

[Install]
WantedBy=multi-user.target
